# perms.list
# chown arg, chmod arg, file list (shell expanded)
0:0 0755 etc/rc.d/init.d/cloud-* var/lib/cloud usr/lib/cloud-init/uncloud-init
x   -x   etc/rc.d/init.d/firewall
